# feeds/luci/luci-app-oled/root/etc/hotplug.d/button/99-oled
#!/bin/sh
# Invoked from cm5-button-scripts /etc/rc.button/* via hotplug-call (procd does not run this file directly).
#
# UCI maps physical buttons to page navigation:
#   menu_select_button — next page (default USERKEY / wps)
#   menu_nav_button    — previous page (default MaskROM / BTN_2)
#
# Future HAT joystick: add GPIO hotplug entries sending up/down/back via oledd-event.sh.

[ "$ACTION" = "pressed" ] || exit 0

nav=$(uci -q get oled.@oled[0].menu_nav_button); [ -z "$nav" ] && nav=BTN_2
sel=$(uci -q get oled.@oled[0].menu_select_button); [ -z "$sel" ] && sel=wps

case "$BUTTON" in
"$sel")
	[ "$sel" != "none" ] && /usr/lib/oled/oledd-event.sh next
	;;
"$nav")
	/usr/lib/oled/oledd-event.sh prev
	;;
esac

exit 0
