# feeds/luci/luci-app-peripherals/root/etc/uci-defaults/90-luci-peripherals-enable
#!/bin/sh
[ -x /etc/init.d/luci-peripherals ] && /etc/init.d/luci-peripherals enable
exit 0
