# feeds/luci/luci-app-peripherals/root/etc/uci-defaults/91-luci-peripherals-drop-crypto
#!/bin/sh
# One-time cleanup after Rockchip crypto UI was removed from luci-app-peripherals.
uci -q delete luci_peripherals.@peripherals[0].rockchip_crypto 2>/dev/null && uci -q commit luci_peripherals
rm -f /etc/modprobe.d/70-luci-peripherals-rockchip-crypto.conf
exit 0
