# feeds/luci/luci-app-oled/root/etc/uci-defaults/oled
#!/bin/sh

. /lib/functions/system.sh
. /usr/lib/oled/cm5-apply-config.sh

# Probe panel after RST; used for diagnostics only — do not gate enable on result.
cm5_oled_probe_panel() {
	local retries=0
	[ -c /dev/i2c-7 ] || return 1
	cm5_oled_release_rst
	while [ "$retries" -lt 5 ]; do
		i2cget -y 7 0x3c 0x00 b >/dev/null 2>&1 && return 0
		retries=$((retries + 1))
		sleep 1
		cm5_oled_release_rst
	done
	return 1
}

uci -q delete oled.@oled[0].chip 2>/dev/null

case "$(board_name)" in
xunlong,orangepi-cm5-base)
	if [ -c /dev/i2c-7 ]; then
		cm5_apply_oled_config /dev/i2c-7
		cm5_oled_release_rst
		cm5_oled_probe_panel || true
	else
		uci -q batch <<-EOF >/dev/null
			set oled.@oled[0].showmenu="1"
			set oled.@oled[0].enable="0"
			set oled.@oled[0].menu_mode="1"
			set oled.@oled[0].path="/dev/i2c-7"
			set oled.@oled[0].ipifname="br-lan"
			set oled.@oled[0].netsource="br-lan"
			set oled.@oled[0].cm5_menu_migrated="1"
			commit oled
EOF
	fi
	;;
esac

if [ "$(uci -q get oled.@oled[0].enable)" = "1" ]; then
	oled_start_enabled
fi

oled_ucitrack_init

sed -e "/\/etc\/init\.d\/oled/d" -i /etc/crontabs/root && /etc/init.d/cron reload
rm -rf /tmp/luci-*
exit 0
